# assets/cache.ini
[program:app]
command=/emojify-cache

[program:consul]
command=/usr/local/bin/consul agent -config-dir=/etc/consul

[program:envoy]
command=/usr/local/bin/consul connect envoy -sidecar-for=emojify-cache
autostart=false

[program:statsd]
command=/bin/statsd_exporter