# assets/consul.ini
[program:consul]
command=/usr/local/bin/consul agent -config-dir=/etc/consul
autostart=false