# assets/ingress.ini
[program:app]
command=/usr/sbin/nginx -g "daemon off";

[program:consul]
command=/usr/local/bin/consul agent -config-dir=/etc/consul

[program:envoy]
command=/usr/local/bin/consul connect envoy -sidecar-for=emojify-ingress
autostart=false

[program:statsd]
command=/bin/statsd_exporter