# assets/ingress.conf
[program:app]
command=/usr/sbin/nginx -g "daemon off";

[program:consul]
command=/usr/local/bin/consul agent -config-dir=/etc/consul