# assets/emojify-api/api.ini
[program:app]
command=/service/emojify-api

[program:consul]
command=/usr/local/bin/consul agent -config-dir=/etc/consul

[program:envoy]
command=/usr/local/bin/consul connect envoy -sidecar-for=emojify-api
autostart=false

[program:statsd]
command=/bin/statsd_exporter